# Library/LaunchAgents/co.lorey.jms-chat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>co.lorey.jms-chat</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>exec /usr/bin/python3 "$HOME/.agents/skills/jms/scripts/jms_chat.py"</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>SOPS_AGE_KEY_FILE</key>
        <string>/Users/william/.config/sops/age/keys.txt</string>
    </dict>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>60</integer>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>/tmp/jms-chat.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/jms-chat.launchd.log</string>
</dict>
</plist>
